# codex_home_assistant/rootfs/etc/s6-overlay/s6-rc.d/codex-ha-init/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -Eeuo pipefail

exec /usr/local/bin/codex-ha-init
