# codex_home_assistant/rootfs/etc/s6-overlay/s6-rc.d/ingress/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
set -Eeuo pipefail

exit_code=${1:-1}
if [[ "${exit_code}" -ne 0 && "${exit_code}" -ne 256 ]]; then
  bashio::log.fatal "Ingress proxy stopped unexpectedly with exit code ${exit_code}"
  printf '%s\n' "${exit_code}" > /run/s6-linux-init-container-results/exitcode
  exec /run/s6/basedir/bin/halt
fi

bashio::log.warning 'Ingress proxy stopped; S6 will restart it'
